# deploy/systemd/recordflow-staging-worker.service
[Unit]
Description=RecordFlow staging worker
Requires=recordflow-staging-migrate.service
Wants=network-online.target
After=network-online.target postgresql.service recordflow-staging-migrate.service

[Service]
Type=simple
User=recordflow-staging
Group=recordflow-staging
WorkingDirectory=/opt/recordflow-staging
EnvironmentFile=/etc/recordflow/staging.env
ExecStart=/opt/recordflow-staging/.venv/bin/python -m recordflow_agent.worker --poll-seconds 1
Restart=always
RestartSec=3
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
